# fly-flash-automator.service
[Unit]
Description=FlyOS Flash Automator - Automatic firmware flashing service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStartPre=/bin/sleep 10
ExecStart=/bin/bash -c "/data/flyos-flash-automator/flash-with-network.sh"
TimeoutStartSec=600
RemainAfterExit=no
User=root
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target