# fly-flash-automator-with-logging.service
[Unit]
Description=FlyOS Flash Automator with Logging
After=multi-user.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStartPre=/bin/sleep 1
ExecStart=/bin/bash -c "/opt/flyos-flash-automator/run-with-logging.sh"
TimeoutStartSec=600
RemainAfterExit=no
User=root
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target